# private_dot_config/systemd/user/dev-git-fetch.service
[Unit]
Description=Fetch and prune active development repositories
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.local/bin/dev

[Service]
Type=oneshot
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:%h/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=SSH_AUTH_SOCK=%h/.1password/agent.sock
ExecStart=%h/.local/bin/dev git-sync --quiet
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
