# private_dot_config/systemd/user/dev-git-fetch.service
[Unit]
Description=Fetch and prune active development repositories
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.local/bin/dev

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dev git-sync --quiet
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
