# private_dot_config/systemd/user/dev-git-fetch.service
[Unit]
Description=Fetch and prune active development repositories
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.local/bin/dev

[Service]
Type=oneshot
Environment=PATH={{ includeTemplate "dev-env-policy.tmpl" . | trim }}
{{- if eq (includeTemplate "onepassword-enabled.tmpl" . | trim) "true" }}
Environment=SSH_AUTH_SOCK={{ includeTemplate "onepassword-agent.tmpl" . | trim }}
{{- end }}
ExecStart=%h/.local/bin/dev git-sync --quiet
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
